# go_exp.service
[Unit]
Description=EXP_GO service

[Service]
ExecStart=/opt/go_exp/main &

